# service/s3sync.service
[Unit]
Description=Sync local folder with S3 via server API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/s3sync/sync.py
WorkingDirectory=/opt/s3sync
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target